# yesterdaypy.timer
[Unit]
Description=YesterdayPy timer

[Timer]
OnCalendar=*-*-* 00:00:00
Unit=yesterdaypy.service
